# deploy/glava-admin.service
[Unit]
Description=GLAVA Admin Panel
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/glava
ExecStart=/opt/glava/venv/bin/gunicorn -w 1 -b 127.0.0.1:5001 admin.app:app
Restart=always
RestartSec=10
EnvironmentFile=/opt/glava/.env

[Install]
WantedBy=multi-user.target
